# modules/metaplugin_ci/files/etc/rc.local
#!/bin/sh -e

bash /opt/metaplugin-ci/files/jenkins-slave.sh online

exit 0
